# deploy/ugassistant.service.template
[Unit]
Description=UGAssistant local assistant
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Group=@USER@
WorkingDirectory=@PROJECT_ROOT@
Environment=PYTHONPATH=@PROJECT_ROOT@/src
ExecStart=@PROJECT_ROOT@/.venv/bin/python -m ugassistant
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
